# deploy/systemd/llm-router-go.service
# /etc/systemd/system/llm-router-go.service
#
# Go router (replacement for litellm-proxy). Binds :4010 — the port the
# legacy LiteLLM proxy used — so client configs don't need to change at
# cutover. To roll back:
#     sudo systemctl stop llm-router-go
#     sudo systemctl start litellm-proxy
#
# Install (typically driven by deploy/scripts/deploy-router.sh):
#     sudo install -m 644 llm-router-go.service \
#         /etc/systemd/system/llm-router-go.service
#     sudo systemctl daemon-reload
#     sudo systemctl enable --now llm-router-go
#
# Requires the existing `llm-router` system user (created for the
# Python agent: `useradd --system --shell /usr/sbin/nologin llm-router`).
# Request logging goes to the dedicated `llm-router-postgres` container
# (postgres:17 on 127.0.0.1:5433) via --postgres-dsn=${ROUTER_PG_DSN}; the
# DSN is materialized into proxy.env by deploy-router-secrets.sh from the
# `ho secret llm-router/reqlog-pg-password`. Empty/unset → NopSink (no log).

[Unit]
Description=LLM Router (Go) — OpenAI-compatible proxy
# docker.service: the reqlog Postgres runs as the `llm-router-postgres`
# container (--restart unless-stopped). Ordering after dockerd shrinks the
# boot race; if the DB still isn't ready the router logs a WARN and runs
# without reqlog (NopSink) — inference never blocks on the log DB.
After=network-online.target docker.service
Wants=network-online.target
# Don't fight with the LiteLLM unit. If it's enabled, this will block
# until it stops; the deploy script disables litellm-proxy first.
Conflicts=litellm-proxy.service

[Service]
Type=simple
User=llm-router
Group=llm-router

ExecStart=/usr/local/bin/llm-router-go \
    --addr=:4010 \
    --models-yaml=/home/erewhon/Projects/erewhon/llm-router/models.yaml \
    --mode=default \
    --log-format=json \
    --log-level=info \
    --postgres-dsn=${ROUTER_PG_DSN} \
    --wellknown-provider-id=llm \
    "--wellknown-provider-name=LLM Router" \
    --wellknown-base-url=https://llm.bcc.sh/v1 \
    --dashboard \
    --dashboard-addr=:4011 \
    --dashboard-public-url=https://llm.bcc.sh

# --dashboard serves the status UI on :4011 — the port the retired Python
# litellm-dashboard used, so the hub Caddy front door (llm-dashboard.bcc.sh
# → euclid's Tailscale IP :4011, gated by oauth2-proxy) needs no change. The
# listener carries NO bearer auth (matching the old Python behaviour); it binds
# all interfaces because Caddy reaches it over Tailscale, and UFW + oauth2-proxy
# are the access boundary. --dashboard-public-url only sets the copy-paste base
# URL shown in the dashboard's Connection card.

# Secrets land here from `ho secret get llm-router/api-key` via
# deploy/scripts/deploy-router-secrets.sh. The `-` prefix makes the unit
# start (with auth disabled, WARN logged) even if the file is missing —
# better than failing closed when the env-file is the wrong shape.
#
#   ROUTER_API_KEYS=sk-...,sk-other-key       # comma-separated; gates /v1/*
#   WELLKNOWN_API_KEY=sk-...                  # surfaced in /.well-known/opencode
#   OPENCODE_ZEN_API_KEY=...                  # upstream key for external Zen models
#   ROUTER_PG_DSN=postgres://router:...@127.0.0.1:5433/router?sslmode=disable
EnvironmentFile=-/etc/llm-router/proxy.env

Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# --- State directories (systemd creates and chowns these) ---
StateDirectory=llm-router-go
RuntimeDirectory=llm-router-go

# --- Filesystem sandboxing ---
# Read-only /usr, /boot, /etc; writes only to StateDirectory/RuntimeDirectory.
ProtectSystem=strict
# Mask /home with a tmpfs, then bind-mount only what the binary needs.
ProtectHome=tmpfs
BindReadOnlyPaths=/home/erewhon/Projects/erewhon/llm-router/models.yaml
PrivateTmp=yes

# --- Kernel protections ---
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes

# --- Capability and privilege restrictions ---
NoNewPrivileges=yes
CapabilityBoundingSet=
AmbientCapabilities=

# --- System call restrictions ---
# NOTE: SystemCallFilter removed 2026-06-06 — combined with the AF_*
# restriction below it broke Go's net resolver on systemd-resolved
# (dial tcp: lookup *: read udp 127.0.0.1:*->127.0.0.53:53: connection refused).
# Reintroduce carefully (test with a real backend call) if hardening matters.
SystemCallArchitectures=native

# --- Namespace and address restrictions ---
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes

# --- Memory ---
LockPersonality=yes
UMask=0077

# --- Network: bound to LAN + Tailscale + NetBird IPs (no internet egress
# needed; backends are LAN-direct, /metrics is scraped locally). Allow
# 100.64.0.0/10 for Tailscale CGNAT space, 100.64.0.0/10 also covers
# NetBird 100.64-100.127 range. ---
# (No IPAddressAllow here — backends live on multiple network ranges
# (192.168.42.0/24 LAN, eventually NetBird mesh). Filtering here would
# break legitimate traffic. UFW on euclid scopes inbound :4010.)

[Install]
WantedBy=multi-user.target
